# chef-repo/cookbooks/transitime/templates/default/core/predictions.service.erb
[Unit]
Description=Predictions Service
After=network.target

[Service]
ExecStart=/var/lib/oba/transitime/core/core.sh
User=ubuntu
